# deploy/happy-paws-bot.service
[Unit]
Description=Happy Paws Grooming Telegram Bot
After=network-online.target happy-paws-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/happy-paws/app
EnvironmentFile=/etc/happy-paws/env
Environment=NODE_OPTIONS=--disable-warning=ExperimentalWarning
ExecStart=/usr/bin/node /opt/happy-paws/app/bot/happy-paws-bot.js
Restart=always
RestartSec=5
User=happy-paws
Group=happy-paws

[Install]
WantedBy=multi-user.target
